# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticpotpy]
directory=/opt/glastopf/
command=bash -c 'exec cd /opt/ElasticpotPY/ && exec python main.py'
autorestart=true

