# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticpotpy]
directory=/opt/ElasticpotPY/
command=/bin/bash -c "python main.py"

