# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticpotpy]
directory=/opt/glastopf/
command=/bin/bash -c "cd /opt/ElasticpotPY/ && python main.py"

