# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticpotpy]
directory=/opt/glastopf/
command=bash -c 'exec /opt/ElasticpotPY/start.sh'
autorestart=true

