# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticpotpy]
directory=/opt/ElasticpotPY/
command=/bin/bash -c "python3 main.py"

